# daemontools/log/run
#!/bin/sh
exec setuidgid spike multilog t ./main
